# container/daycare/build/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=warn
user=root

[program:caldav_sync]
command=python -u /app/caldav_sync.py
user=daycare
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=CONFIG_PATH="/config/config.yml",MESSAGES_PATH="/config/messages.yml",STATE_PATH="/data/state.json",IDEAS_PATH="/data/ideas.json"

[program:notifier]
command=python -u /app/notifier.py
user=daycare
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=CONFIG_PATH="/config/config.yml",MESSAGES_PATH="/config/messages.yml",STATE_PATH="/data/state.json",IDEAS_PATH="/data/ideas.json",LOG_LEVEL="INFO"

[program:receiver]
command=python -u /app/receiver.py
user=daycare
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=CONFIG_PATH="/config/config.yml",MESSAGES_PATH="/config/messages.yml",STATE_PATH="/data/state.json",IDEAS_PATH="/data/ideas.json"

[program:cleanup_runner]
command=python -u /app/cleanup_runner.py
user=daycare
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=CONFIG_PATH="/config/config.yml",MESSAGES_PATH="/config/messages.yml",STATE_PATH="/data/state.json",IDEAS_PATH="/data/ideas.json",ARCHIVE_DIR="/data"
